# infra/ansible/templates/openclaw/openclaw-gateway.service.j2
[Unit]
Description=OpenClaw Gateway ({{ openclaw_agent }})
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ openclaw_user }}
Group={{ openclaw_user }}
WorkingDirectory={{ openclaw_home }}
Environment=HOME={{ openclaw_home }}
Environment=NODE_ENV=production
Environment=ANTHROPIC_API_KEY={{ anthropic_api_key }}
Environment=GEMINI_API_KEY={{ gemini_api_key }}
Environment=GOG_KEYRING_PASSWORD=openclaw-grimes
Environment=GOG_ACCOUNT=smcgroarty96@gmail.com
ExecStart=/usr/bin/openclaw gateway --port 18789
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths={{ openclaw_home }} /usr/local/bin /tmp
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
